# testing/ndppd/ndppd.init
#!/sbin/openrc-run

name="$SVCNAME"
command="/usr/sbin/$SVCNAME"
pidfile="/var/run/$SVCNAME.pid"
command_args="-p $pidfile -d"

depend() {
	need net
}

